# app.service
[Unit]
Description=HolidayHub Website
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/HolidayHub
ExecStart=/bin/bash -c 'source env/bin/activate && source .env.production && gunicorn -b 0.0.0.0:8000 server:create_app()'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target